# .docker/supervisor/confs/conf.d/program_clock3.conf
[program:clock3]
command=/opt/cesi-dev/bin/clock.sh no3
autostart=true
autorestart=true
stderr_logfile=/opt/cesi-dev/logs/clock3.err.log
stdout_logfile=/opt/cesi-dev/logs/clock3.out.log
