# .docker/supervisor/confs/conf.d/program_clock2.conf
[program:clock2]
command=/opt/cesi-dev/bin/clock.sh no2
autostart=true
autorestart=true
stderr_logfile=/opt/cesi-dev/logs/clock2.err.log
stdout_logfile=/opt/cesi-dev/logs/clock2.out.log